# scripts/deploy_init/idp_shop.service
[Unit]
Description=idp_shop

[Service]
ExecStart=/usr/local/bin/idp_shop -c /etc/idp_shop.toml

[Install]
WantedBy=multi-user.target
